# roles/monitoring/node_exporter/templates/run-textfile-collectors@.service
[Service]
Type=simple
ExecStart={{ prom_libexec_dir }}/run-textfile-collectors.sh {{ prom_libexec_dir }}/textfile_collectors/%i
