# roles/monitoring/node_exporter/files/run-textfile-collectors-1min.timer
[Timer]
OnCalendar=minutely
Unit=run-textfile-collectors@1min.service

[Install]
WantedBy=timers.target
