# roles/monitoring/node_exporter/files/run-textfile-collectors-daily.timer
[Timer]
OnCalendar=daily
Unit=run-textfile-collectors@daily.service

[Install]
WantedBy=timers.target
