# roles/monitoring/node_exporter/files/run-textfile-collectors-hourly.timer
[Timer]
OnCalendar=hourly
Unit=run-textfile-collectors@hourly.service

[Install]
WantedBy=timers.target
